# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0

[program:backend]
command=dotnet /app/backend/VManBackend.dll
autostart=true
autorestart=true
# Lower priority number = starts first; backend runs migrations before accepting traffic
priority=100
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:frontend]
command=node /app/frontend/server.js
autostart=true
autorestart=true
priority=200
environment=PORT=3000,HOSTNAME=127.0.0.1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
# Wait for the backend /health endpoint before accepting traffic, so nginx doesn't
# serve 502s for /api/* while migrations are running on startup.
command=sh -c "until curl -fsS http://127.0.0.1:8080/health > /dev/null 2>&1; do sleep 2; done; exec nginx -g 'daemon off;'"
autostart=true
autorestart=true
priority=300
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
